# deploy/rep-update.service
[Unit]
Description=Ryanne REP — Atualização automática
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rep-ctl update
StandardOutput=journal
StandardError=journal
